# .vim/sonictemplate/upstart/base-service.conf
description "{{_name_}}service"
author  "root <root@localhost>"

start on runlevel [2345]
stop on runlevel [016]

# chdir /root/workdir
# exec {{_name_}} >> /var/log/{{_name_}}.log 2>&1
respawn

